# .config/.config/systemd/user/battery-check.timer
[Unit]
Description=Periodic battery check

[Timer]
OnBootSec=30s
OnUnitActiveSec=300s

[Install]
WantedBy=timers.target
