# .config/.config/systemd/user/battery-check.service
[Unit]
Description=Battery check

[Service]
Type=oneshot
ExecStart=%h/common/repos/dotfiles/scripts/battery_test.sh
